# configs/system/Fedora/profile.d/which2.sh
# Initialization script for bash and sh

# export AFS if you are in AFS environment
alias which='(alias; declare -f) | /usr/bin/which --tty-only --read-alias --read-functions --show-tilde --show-dot'
